# rootfs/etc/s6-overlay/s6-rc.d/generate-prom/run
#!/bin/sh
exec /scripts/generate-prom.sh
